# etc/profile.d/sys-update.sh
alias sysupdate="rpm-ostree update && flatpak update -y && distrobox upgrade --all"